# docs/alacritty-circadian.service
[Unit]
Description=Alacritty automatic theme switch

[Service]
ExecStart=$HOME/.local/bin/alacritty-circadian.py

[Install]
WantedBy=default.target
